# ansible/roles/pi-minimised/files/hdmi-off.service
[Unit]
Description=Disable HDMI
DefaultDependencies=no

[Service]
Type=simple
ExecStart=/usr/bin/tvservice -o

[Install]
WantedBy=sysinit.target